# sentinel/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
user=root
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:master]
command=redis-server $MASTER_CONFIG_FILE

[program:slave1]
command=redis-server $SLAVE1_CONFIG_FILE

[program:slave2]
command=redis-server $SLAVE2_CONFIG_FILE

[program:sentinel1]
command=redis-server $SENTINEL1_CONFIG_FILE --sentinel

[program:sentinel2]
command=redis-server $SENTINEL2_CONFIG_FILE --sentinel

[program:sentinel3]
command=redis-server $SENTINEL3_CONFIG_FILE --sentinel
